# synapse_nerve_zookeeper/configuration/lb.conf
description "synapse server process"

start on (net-device-up IFACE=eth0 or runlevel [2345])
stop on runlevel [!12345]

respawn

pre-start script
  logger "pre-start for synapse"
end script

post-start script
  logger "post-start for synapse"
end script

script
  exec su -s /bin/bash -c '/usr/local/rvm/gems/ruby-2.2.1@synapse/wrappers/synapse -c /vagrant/configuration/$(hostname).conf.json' vagrant
end script